# systemd/geoip-block.timer
[Unit]
Description=Run geoip-block periodically

[Timer]
OnCalendar=daily
Persistent=true

[WantedBy]
WantedBy=timers.target
